# haoStock_deb/DEBIAN/prerm
#!/bin/bash

systemctl disable haoStock.timer
systemctl stop haoStock.timer
